# vt.archlinux/docker/templates/docker.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// -s {{ docker_storage_driver }}
